# bootstrap/agent-factory-reporter.service
[Unit]
Description=Factory AI hourly progress report
After=agent-factory-worker.service

[Service]
Type=oneshot
SyslogIdentifier=factory-ai-reporter
User=factory
Group=factory
WorkingDirectory=/opt/agent-factory/app
ExecStart=/usr/bin/node /opt/agent-factory/app/src/run-with-env.js /etc/agent-factory.env /usr/bin/node /opt/agent-factory/app/src/reporter.js
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agent-factory/state
UMask=0027
