# bootstrap/factory-ai-snapshot.service
[Unit]
Description=Publish Factory AI operator dashboard snapshot
After=network-online.target agent-factory-control.service
Wants=network-online.target

[Service]
Type=oneshot
User=factory
Group=factory
SupplementaryGroups=systemd-journal
WorkingDirectory=/opt/agent-factory/app
ExecStart=/usr/bin/node /opt/agent-factory/app/src/run-with-env.js /etc/agent-factory-control.env /usr/bin/node /opt/agent-factory/app/src/snapshot.js
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
UMask=0027
SyslogIdentifier=factory-ai-snapshot
